# supervisord/php-fpm.conf
[program:php-fpm]
command=/usr/sbin/php5-fpm --nodaemonize --fpm-config /etc/php5/fpm/php-fpm.conf
numprocs=1
autostart=true
autorestart=true