# supervisord/php-fpm.conf
[program:php-fpm]
command=/usr/sbin/php-fpm7.4 --nodaemonize --fpm-config /etc/php/7.4/fpm/php-fpm.conf
numprocs=1
autostart=true
autorestart=true